# Settings/crontab
0 20 * * * PATH/Backups/dexters_cloud.sh
0 20 * * 2 PATH/Backups/settings_backup.sh

# last line should be blank for cronetab 
